# image-builder/profiles/fallback/rc.local
mkdir -p /tmp/apt/lists
mkdir -p /tmp/apt/cache

/usr/local/bin/smarties_update

exit 0
