# debian/ohos.cron.d
#
# Regular cron jobs for the ohos package
#
*/5 * * * * root /usr/bin/ohos-auto

