# image-builder/profiles/skel/rootfs/etc/rc.local
/usr/local/bin/clear_update_flag

exit 0

